# deploy/systemd/portfolio-bot.service
[Unit]
Description=Portfolio Telegram Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/portfolio-telegram-analytics
EnvironmentFile=/opt/portfolio-telegram-analytics/.env

# Prevent duplicates (optional but recommended)
ExecStart=/usr/bin/flock -n /var/lock/portfolio-bot.lock /opt/portfolio-telegram-analytics/.venv/bin/python -m portfolio_app.telegram_bot

Restart=on-failure
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=multi-user.target