# deploy/systemd/portfolio-weekly.timer
[Unit]
Description=Run portfolio weekly summaries on Sundays

[Timer]
OnCalendar=Sun *-*-* 20:00:00
Persistent=true

[Install]
WantedBy=timers.target